# skiffos_project/my_docker_config/root_overlay/etc/systemd/system/mnt-hdd_usb.mount
[Unit]
Description=mounting Files from SSD

[Mount]
Where=/mnt/sdb1
What=/dev/disk/by-label/storage1

[Install]
WantedBy=multi-user.target